# systemd/asus-numpad@.service.in
[Unit]
Description=Virtual numpad driver for %I

[Service]
Type=simple
ExecStart=@BIN@ %I
